# modules/packages/templates/election-orchestra.conf.erb
[program:eorchestra]
command=/home/eorchestra/venv/bin/uwsgi --ini /home/eorchestra/election-orchestra/auth.ini
directory=/home/eorchestra/election-orchestra/
autostart=true
autorestart=true
user=eorchestra
environment=PATH="/home/eorchestra/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
redirect_stderr=true
environment='FRESTQ_SETTINGS'="base_settings.py"